# systemd/stadia-controller.service
[Unit]
Description=Stadia Controller Auto-Connect Service
Documentation=man:bluetoothctl(1)
After=bluetooth.service
Wants=bluetooth.service
BindsTo=bluetooth.service
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/andrewmeckley/ros2_ws/src/rover1/scripts/connect_stadia_controller.sh
User=andrewmeckley
Group=andrewmeckley

# Environment
Environment=PATH=/usr/bin:/bin:/usr/sbin:/sbin

# Timeout settings
TimeoutStartSec=120
TimeoutStopSec=30

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stadia-controller

[Install]
WantedBy=multi-user.target